# scripts/com.investment-analytics.watchlist.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS launchd job — runs the watchlist snapshot every weekday at
  21:30 IST (16:00 UTC). 21:30 IST is after AMFI publishes EOD NAV,
  so the snapshot reflects the day's actual closing values.

  To install:
    cp scripts/com.investment-analytics.watchlist.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.investment-analytics.watchlist.plist

  To uninstall:
    launchctl unload ~/Library/LaunchAgents/com.investment-analytics.watchlist.plist
    rm ~/Library/LaunchAgents/com.investment-analytics.watchlist.plist

  Edit the WorkingDirectory + Program Arguments below to match your
  install path before installing.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.investment-analytics.watchlist</string>

    <key>WorkingDirectory</key>
    <string>/Users/snapdevio/Documents/investment-analytics</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/snapdevio/Documents/investment-analytics/venv/bin/python</string>
        <string>-m</string>
        <string>backend.jobs.watchlist</string>
        <string>run</string>
    </array>

    <!--
      Run weekdays (Mon-Fri) at 21:30 local time. AMFI publishes EOD
      NAV between 22:00 and 23:00 IST so this scheduling assumes IST
      machine time; if your laptop's clock is in another timezone,
      adjust Hour accordingly.
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/snapdevio/Documents/investment-analytics/logs/watchlist.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/snapdevio/Documents/investment-analytics/logs/watchlist.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
